# Halloween-Visions-Xavier/halloween.service
[Unit]
Description=Halloween Hand Detection Projection (Denhac)
After=graphical.target
Wants=graphical.target
# Wait for display server to be ready
After=display-manager.service

[Service]
Type=simple
User=colin
Group=colin

# Environment variables for X11 display
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/colin/.Xauthority"

# Working directory
WorkingDirectory=/home/colin/Documents/MLVisionsProjects/Halloween-Visions-Projection

# Wait 15 seconds after graphical target to ensure display is ready
ExecStartPre=/bin/sleep 15

# Main command - using venv python with validated parameters
ExecStart=/home/colin/Documents/MLVisionsProjects/Halloween-Visions-Projection/venv/bin/python3 simple_projection.py --model models/Colin1.pt --video-sleep videos/sleeping_face.mp4 --video-scare videos/angry_face.mp4 --conf 0.7 --source 0 --fullscreen

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=halloween-projection

[Install]
WantedBy=graphical.target
